# raycast/context-dictionary/scripts/com.1206yaya.raycast-context-dictionary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.1206yaya.raycast-context-dictionary</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/zak/ghq/github.com/1206yaya/dotfiles/raycast/context-dictionary/scripts/start-dev.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/raycast-context-dictionary.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/raycast-context-dictionary.log</string>
</dict>
</plist>
